# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:crond]
command = /usr/sbin/cron
user = root
autostart = true