# supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command = /usr/sbin/cron
user = root
autostart = true